# input-20260323-133834/drop-ins/etc/systemd/system/httpd.service.d/override.conf
[Service]
TimeoutStartSec=600
LimitNOFILE=65535
