# input-20260323-133834/config/etc/cron.daily/cleanup.sh
#!/bin/sh
# Cleanup script — driftify synthetic fixture
find /tmp -name 'myapp-*' -mtime +7 -delete
