# input-20260323-133834/drop-ins/etc/systemd/system/nginx.service.d/override.conf
[Service]
LimitNOFILE=131072
ExecStartPost=/usr/local/bin/notify-deploy.sh
