# input-20260323-133834/config/etc/cron.d/backup-daily
# Daily backup — driftify synthetic fixture
# yoinkc should convert this to a systemd timer
0 2 * * * root /usr/local/bin/backup.sh >> /var/log/myapp/backup.log 2>&1
