# input-20260323-133834/config/etc/cron.d/complex-job
# Complex cron with env deps — driftify kitchen-sink fixture
MAILTO=ops@example.com
APP_ENV=production
30 6 * * 1-5 appuser /opt/myapp/scripts/weekday-report.sh
