# tests/fixtures/host_etc/usr/lib/systemd/system/test-static.service
[Unit]
Description=Test static service (no Install section)

[Service]
ExecStart=/bin/true
