# tests/fixtures/host_etc/usr/lib/systemd/system/test-installable.service
[Unit]
Description=Test installable service

[Service]
ExecStart=/bin/true

[Install]
WantedBy=multi-user.target
