# tests/fixtures/host_etc/usr/lib/systemd/system/test-masked.service
[Unit]
Description=Test service that will be masked

[Service]
ExecStart=/bin/true

[Install]
WantedBy=multi-user.target
