# input-20260323-133834/config/etc/systemd/system/cron-daily-cleanup-sh.timer
[Unit]
Description=Generated from cron.daily: etc/cron.daily/cleanup.sh
# Original: cron.daily script

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
