# input-20260323-133834/config/etc/systemd/system/cron-backup-daily.timer
[Unit]
Description=Generated from cron: etc/cron.d/backup-daily
# Original cron: 0 2 * * *

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
