# input-20260323-133834/config/etc/systemd/system/cron-appuser.timer
[Unit]
Description=Generated from cron: var/spool/cron/appuser
# Original cron: */15 * * * *

[Timer]
OnCalendar=*-*-* *:*/15:00
Persistent=true

[Install]
WantedBy=timers.target
