# input-20260323-133834/config/etc/systemd/system/cron-daily-cleanup-sh.service
[Unit]
Description=Timer from cron.daily etc/cron.daily/cleanup.sh

[Service]
Type=oneshot
ExecStart=/etc/cron.daily/cleanup.sh
