# input-20260323-133834/config/etc/systemd/system/driftify-backup.service
[Unit]
Description=Run backup script

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup.sh
