# input-20260323-133834/config/etc/systemd/system/myapp-report.service
[Unit]
Description=MyApp Daily Report
After=network.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/generate-report.sh
StandardOutput=journal
